# packaging/apt/debian_server/debian/preinst
#!/bin/sh

set -e

if [ -f /etc/init.d/scatterbytes-server ]; then
    invoke-rc.d scatterbytes-server stop
fi
